# services/discord-bot/chimera/chimera.container
[Unit]
Description=Chimera Discord Bot
After=network-online.target

[Container]
ContainerName=chimera
Image=ghcr.io/baguettejet/chimera:latest
AutoUpdate=registry

WorkingDir=/app

Volume=%h/discord-bots/chimera/repo:/app:Z
EnvironmentFile=%h/discord-bots/chimera/repo/.env

# fix terminal colour for Cockpit and run main script
Exec=/usr/bin/script -q -c "python -u bot.py" /dev/null

[Service]
Restart=on-failure

[Install]
WantedBy=default.target